# deploy/systemd/last-mile-studio-certbot-renew.service
[Unit]
Description=Renew the AfterPrompt IP certificate
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/certbot/bin/certbot renew --quiet
PrivateTmp=true
NoNewPrivileges=true
